# Services/source-whatsapp.service
[Unit]
Description=Listens and forwards WhatsApp Messages
After=network.target
OnFailure=unit-failed@%n
#ConditionPathExists=!/etc/ssh/sshd_not_to_be_run

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m smsgateway.sources.whatsapp
User=smsd
Restart=always
RestartSec=20
StartLimitInterval=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
